# dist/dnf-notify.service
[Unit]
Description=Check and notify about dnf updates
Wants=dnf-notify.timer
Requires=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/dnf-notify
